# supervisor/queue_log.ini
[program:queue_log]                                                                                                                                                                                
directory=/opt/asterisk14_sql_queue_log
command=python /opt/asterisk14_sql_queue_log/queue_log2sql_insert.py /var/log/asterisk/queue_log 'postgres://postgres@/asterisk'
user=root
autorestart=true
