# supervisor/queue_log.ini
[program:queue_log]                                                                                                                                                                                
directory=/opt/asterisk14_sql_queue_log
command=/opt/asterisk14_sql_queue_log/queue_log.sh
user=asterisk
autorestart=true
