# tornado/queue_log.ini
[program:queue_log]                                                                                                                                                                                
directory=/home/andriy
command=/home/andriy/queue_log.sh
user=root 
autorestart=true
 